# source/testing/SOURCES/dashcore-0.13.0-contrib/linux/systemd/usr-lib-systemd-system_dashd.service
[Unit]
Description=Dash Core Daemon
After=syslog.target network.target

# Notes:
# Do not edit this file. Instead edit /etc/sysconfig/dashd
# ...unless you ever duplicate this service (run multiple-nodes on one
# machine). In that case you need to rename this and deal with the pid file
# changes
#
# Watch the daemon service actions in the syslog journal with:
# sudo journalctl -u dashd.service -f

[Service]
Type=forking
User=dashcore
Group=dashcore
EnvironmentFile=/etc/sysconfig/dashd

#PIDFile=/run/dashd/dashd.pid # PIDFile= is unneccessary (systemd figures it out) and throws false errors in the journal
ExecStart=/usr/bin/dashd -daemon -conf=${CONFIG_FILE} -datadir=${DATA_DIR} -pid=/run/dashd/dashd.pid $OPTIONS
# Experimental: sleep then mnsync reset
#ExecStartPost=/usr/bin/sleep $START_SLEEP_TIME
#ExecStartPost=/usr/bin/dash-cli -conf=${CONFIG_FILE} -datadir=${DATA_DIR} mnsync reset
# Send email after start
ExecStartPost=-/etc/sysconfig/dashd-scripts/dashd.send-email.sh start $EMAIL_FROM $EMAIL_TO $MASTERNODE_ALIAS
# Time that systemd gives a process to start before shooting it in the head
TimeoutStartSec=10m

# If ExecStop is not set, systemd sends a SIGTERM, which is "okay", just not ideal
ExecStop=/usr/bin/dash-cli -conf=${CONFIG_FILE} -datadir=${DATA_DIR} stop
# Need to allow dashd to shutdown properly before taking additional actions.
ExecStopPost=-/usr/bin/echo "Sleeping during shutdown for ${STOP_SLEEP_TIME} seconds"
ExecStopPost=-/usr/bin/sleep $STOP_SLEEP_TIME
# Send email at stop
ExecStopPost=-/etc/sysconfig/dashd-scripts/dashd.send-email.sh stop $EMAIL_FROM $EMAIL_TO $MASTERNODE_ALIAS
# Time that systemd gives a process to stop before shooting it in the head
TimeoutStopSec=120

Restart=on-failure
# If something triggers an auto-restart, let's wait a bit before taking further action
# Note: This value is in addition to the stop sleep time
RestartSec=15

# In this interval span of time, we allow systemd to start dashd "burst" number
# of times. With Dash we really only want one instance started, so... let's
# really limit this. But we want to give systemd some room to attempt to
# correct things. To be honest, I think the way things are configured between
# these settings and TimeoutStartSec, only one instance will be initiated.
StartLimitInterval=300
StartLimitBurst=3

PrivateTmp=true

[Install]
WantedBy=multi-user.target

# Really useful:
# * https://www.digitalocean.com/community/tutorials/understanding-systemd-units-and-unit-files
# * https://www.freedesktop.org/software/systemd/man/systemd.service.html
# * man systemd, man systemd.service, and man systemd.unit
